# roles/docker-registry/files/docker-registry.service
[Service]
ExecStart=/bin/docker-compose -f /opt/my_docker/docker-compose.yml up
Restart=on-failure